# inittab
dnsbackup:23456:wait:/home/shuai.w//.sangfor/ssl/bin/dnsrecover
dnsbackup:23456:wait:/home/shuai.w/.sangfor/ssl/bin/dnsrecover
